# rc.d/chyves
#!/bin/sh

# Add the following lines to /etc/rc.conf to enable chyves:
#
# chyves_enable="YES"
# chyves_flags="pool=zroot kmod=1 net=em0"
#
# PROVIDE: chyves
# REQUIRE: LOGIN sshd
#

. /etc/rc.subr

name="chyves"
rcvar=chyves_enable

# read configuration and set defaults
load_rc_config "$name"
: ${chyves_enable="NO"}
: ${chyves_flags=""}

start_cmd="chyves_start"
stop_cmd="chyves_stop"

chyves_start()
{
	echo "Starting chyves guests..."
	/usr/local/sbin/chyves setup ${chyves_flags}
	sleep 8 # Sleeping for network to come up.
	guestlist="$( zfs get -H -t filesystem -o name,value chyves:rcboot | grep chyves/guests/ | sort -k2nr -k1 | awk '{ if ($2 > 0) print $1; }' )"
	for i in $guestlist ; do
		/usr/local/sbin/chyves start $i
		sleep 3
	done
}

chyves_stop()
{
	echo "Stopping all chyves guests..."
	/usr/local/sbin/chyves scram
}

run_rc_command "$1"
